# overlay/etc/cron.hourly/moodle
#!/bin/bash

su www-data -c "/usr/bin/php /var/www/moodle/admin/cli/cron.php > /dev/null 2>&1"
